# 30-ansible/roles/yatta/templates/etc/profile.d/yatta.sh
YATTA_HOME={{dev_dir}}/yatta
PATH=${YATTA_HOME}/bin:${PATH}

export YATTA_HOME PATH
