# 30-ansible/roles/mvn/templates/etc/profile.d/mvn.sh
M2_HOME={{dev_dir}}/mvn
PATH=${M2_HOME}/bin:${PATH}

export M2_HOME PATH
