# 30-ansible/roles/jdk/templates/etc/profile.d/jdk.sh
JAVA_HOME={{dev_dir}}/java
PATH=${JAVA_HOME}/bin:${PATH}

export JAVA_HOME PATH
